# sphinxsearch/crontab
* * * * * root . /etc/environment; echo "Indexed $(date)" >> /var/log/cron.log 2>&1 && /usr/local/bin/indexer -c /etc/sphinxsearch/sphinx.conf --all --rotate >> /var/log/cron.log 2>&1

